# afharchive.service
[Unit]
Description=AFHArchive Web Application
After=network.target

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/path/to/AFHArchive
Environment=PATH=/path/to/AFHArchive/venv/bin
EnvironmentFile=/path/to/AFHArchive/.env
ExecStart=/path/to/AFHArchive/venv/bin/python run.py production
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
